# services/dashboard.service
[Unit]
Description=Dashboard backend

[Service]
Restart=always
User=root
ExecStart=/bin/dashboard

[Install]
WantedBy=multi-user.target